# build/cron.d/signal
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
PULSE_SERVER=unix:/tmp/pulse/native
TZ=Asia/Tokyo

59 * * * * root sudo -Eu player /app/play.sh Xenoblade2/15_Gormott_Province.flac 1>>/proc/1/fd/1 2>>/proc/1/fd/2
* * * * * root date 1>>/proc/1/fd/1 2>>/proc/1/fd/2
